# service-files/mongo-backup.timer
[Unit]
Description=Runs the mongo backup job

[Timer]
OnCalendar=03:00

[Install]
WantedBy=timers.target

[X-Fleet]
MachineOf=mongodb@2.service
